# build_files/sysfiles/etc/cron.daily/bootc_upgrade.sh
#!/bin/bash
upgrading=$(bootc upgrade 2>/dev/null)
